# meta-common/recipes-intel/sel-logger/files/xyz.openbmc_project.Logging.IPMI.service
[Unit]
Description=IPMI SEL Logging Service

[Service]
Restart=always
ExecStart={bindir}/sel-logger
Type=simple

[Install]
WantedBy=multi-user.target
